# cuckoo-rooter.service
[Unit]
Description=Cuckoo Sandbox Root Socket

[Service]
Type=simple
ExecStart=/usr/local/bin/cuckoo rooter

[Install]
WantedBy=multi-user.target

